# ansible/roles/deploy/templates/wtmm.target
[Unit]
Wants={{ application_name }}_server.service

[Install]
WantedBy=default.target
